# agent/mete-watch.service
[Unit]
Description=mete-watch

[Service]
ExecStart=/usr/bin/python3 /root/watch-mete/watch.py
Restart=always

[Install]
WantedBy=multi-user.target
